# linux/etc/init.d/kiosk
#!/sbin/openrc-run

name="kiosk"
command="/usr/bin/chromium"
command_args="--no-sandbox --kiosk --no-first-run --disable-infobars --noerrdialogs --use-fake-ui-for-media-stream --remote-debugging-port=9222 --enable-features=WebHID --disable-features=WebHIDBlocklist --auto-accept-this-tab-capture --enable-experimental-web-platform-features \"https://hardwave.tgr.rs/index.html?api_key=&chat_id=\""
command_background=true
pidfile="/run/kiosk.pid"
command_user="kiosk:kiosk"
directory="/home/kiosk/harddver"
supervisor="supervise-daemon"
supervise_daemon_args="--respawn-delay 1 --respawn-max 0 --respawn-period 1 --stdout /home/kiosk/logs/kiosk.log --stderr /home/kiosk/logs/kiosk.err"

start_pre() {
    export DISPLAY=:0
}
